# .docker/cron/crontab
# [Minute] [hour] [Day_of_the_Month] [Month_of_the_Year] [Day_of_the_Week] [command]
30 3 * * *   docker exec -it php-qit sh -c "php /var/www/webapp/bin/console app:mail:send"
